# packer/files/gh-runner.service
[Unit]
Description=GitHub Actions ephemeral runner (one job, then poweroff)
Documentation=https://github.com/jeffresc/actions-runner-scaleset-proxmox
# Wait for the network so the runner can reach github.com.
After=network-online.target
Wants=network-online.target
# The path-unit triggers us; we don't auto-start at boot.
ConditionPathExists=/opt/actions-runner/jitconfig.env

[Service]
Type=simple
User=runner
Group=runner
WorkingDirectory=/opt/actions-runner

# The orchestrator writes the JIT config as a systemd env file
# (single `JIT_CONFIG=…` line). Using EnvironmentFile avoids the
# `bash -c '$(cat …)'` shell-substitution path that, for reasons we
# never fully tracked down under Proxmox guest-agent + path-unit
# triggering, produced empty arguments to run.sh under systemd while
# the same command worked perfectly from a shell.
EnvironmentFile=/opt/actions-runner/jitconfig.env

ExecStart=/opt/actions-runner/run.sh --jitconfig ${JIT_CONFIG}

# Regardless of the runner's exit code, power the VM off so the
# orchestrator's destroy path doesn't have to wait on shutdown.
# The `+` prefix tells systemd to run this regardless of the User= above
# (poweroff requires root privileges; the runner user can't issue it).
ExecStopPost=+/usr/bin/systemctl poweroff

Restart=no
TimeoutStartSec=0

# Hardening — these are conservative enough that most CI workloads still
# work. Loosen on a per-deployment basis if a specific workflow needs
# more (e.g. mounting block devices, modifying kernel modules).
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/actions-runner /home/runner /tmp /var/tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectControlGroups=true
# docker/podman in CI need access to namespaces.
RestrictNamespaces=false
# Some build tools legitimately use setuid/setgid binaries.
RestrictSUIDSGID=false
LockPersonality=true
SystemCallArchitectures=native

[Install]
# Triggered by the path unit, not WantedBy=multi-user.target.
